# files/apt-daily.timer
[Timer]
OnCalendar=
OnCalendar=Tuesday 4:30
RandomizedDelaySec=10m
