# files/apt-daily-upgrade.timer
[Timer]
OnCalendar=
OnCalendar=Tuesday 5:00
RandomizedDelaySec=10m
